# docs/itm-r2-sitemaps.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Regenerate a complete static snapshot every six hours. flock prevents overlap;
# sitemap-index.xml switches only after every new shard has uploaded successfully.
15 */6 * * * root if [ -f /opt/indiantrademart-backend/.env.r2-sitemap ]; then flock -n /var/lock/itm-r2-sitemaps.lock bash -lc 'set -a; source /opt/indiantrademart-backend/.env.r2-sitemap; set +a; unset SITEMAP_R2_MAX_URLS SITEMAP_R2_SNAPSHOT_ID SITEMAP_R2_DRY_RUN DRY_RUN SITEMAP_LASTMOD_DATE; export SITEMAP_CRAWLER_BASE_URL=https://indiantrademart.com SITEMAP_CATEGORY_SCOPE=all SITEMAP_R2_URL_LIMIT=50000 SITEMAP_R2_GZIP_LEVEL=6 SITEMAP_R2_RETAIN_SNAPSHOTS=2; cd /opt/indiantrademart-backend && nice -n 10 npm run seo:sitemap:r2:publish' >> /var/log/itm-r2-sitemaps.log 2>&1; fi
